# deploy/cortana-brain.service
# CORTANA Brain — Python: wake word, STT, grammar, gazetteer, incident engine,
# routing, TTS synthesis. Binds the IPC socket; Ears connects to it (GDD §15).
#
# Install: /etc/systemd/system/cortana-brain.service (deploy/install.sh does this).
# Code runs from the /opt/cortana/current symlink, which install.sh flips
# atomically between staged releases (GDD §17) — a restart is all it takes to
# pick up a new release, and --rollback is a flip plus this same restart.
# Token comes ONLY via LoadCredential= (GDD §18/§22, CLAUDE.md constraint 12);
# the process reads $CREDENTIALS_DIRECTORY/token.

[Unit]
Description=CORTANA Brain (fleet intel engine)
Documentation=file:/opt/cortana/current/brain/docs/GDD.md
After=network-online.target
Wants=network-online.target
# A genuinely broken brain must stop and alert, not crash-loop invisibly
# forever: five failed starts inside five minutes stops the unit, which
# lands it in failed state and fires OnFailure=.
StartLimitIntervalSec=300
StartLimitBurst=5
# On failure (start limit hit, config-error exit, crash the restart policy
# refuses) run the alert hook: appends to /var/lib/cortana/alerts.log and
# posts the journal tail to the ops webhook when /etc/cortana/alert_webhook
# exists. The alarm that works precisely when the bot cannot start.
OnFailure=cortana-alert@%n.service

[Service]
Type=simple
# --- Type=notify, ready when the brain actually is -----------------------
# The brain does not call sd_notify yet (that lands with the supervision
# work in brain/). Once it sends READY=1 only after gateway-ready + IPC
# socket bind + engine load, and the health loop pings WATCHDOG=1, replace
# Type=simple above with:
#   Type=notify-reload
#   NotifyAccess=main
#   WatchdogSec=90
# -------------------------------------------------------------------------
User=aura
Group=aura
# Hugging Face cache shared with install.sh's whisper prefetch: a size-name
# stt.model (e.g. "small") resolves through this cache, so the weights staged
# at deploy time are the ones the service loads — never a surprise ~460 MB
# download as aura at the first utterance. Dev boxes without this env fall
# back to ~/.cache/huggingface as usual.
Environment=HF_HOME=/var/lib/cortana/hf
ExecStart=/opt/cortana/current/venv/bin/python -m cortana --config /etc/cortana/cortana.yaml
# SIGHUP → ConfigHolder.reload() (hot-reload cortana.yaml, GDD §16/§18)
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
# Exit 78 (EX_CONFIG: invalid config/gazetteer/routing) and 69
# (EX_UNAVAILABLE: auth) are operator problems a restart cannot fix — stop
# cleanly with the reason in the journal instead of looping every 5s.
RestartPreventExitStatus=78 69
# Brain bounds its own graceful shutdown (~8s); this is the backstop so a wedged
# close still restarts promptly instead of waiting the 90s default before SIGKILL.
TimeoutStopSec=15

# Secrets: token file 0600, surfaced read-only at $CREDENTIALS_DIRECTORY/token.
LoadCredential=token:/etc/cortana/token
# Anthropic API key for the chat override channel (constraint 12). Always
# loaded: install.sh guarantees the source file exists (empty when the
# channel is unused — an empty credential reads as "no key" and the channel
# stays off). Kept unconditional so a re-run of install.sh can never
# silently disable a configured override channel, which is exactly what
# happened when this line shipped commented out and an install.sh re-run
# clobbered a hand-edited unit.
LoadCredential=anthropic:/etc/cortana/anthropic

# Hardening (GDD §22)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/cortana /run/cortana

[Install]
WantedBy=multi-user.target
